# systemd/lexa-modbus.service
[Unit]
# StartLimit* MUST live in [Unit] (systemd >=v230 ignores them in
# [Service]); the default 5-in-10s otherwise applies and a crash
# burst (e.g. the power-cut-retained scenario) permanently disables
# the service. (V1RC FINDING A, 2026-07-06.)
StartLimitIntervalSec=60
StartLimitBurst=5
Description=LEXA Modbus/SunSpec Device Service
Documentation=https://github.com/lexa/lexa-hub
After=network.target mosquitto.service
# Wants= not Requires=: a clean mosquitto stop/restart must NOT stop-propagate
# to the lexa services — they survive broker restarts via paho auto-reconnect
# (the crash-only/reconnect design; mqtt-broker-restart scenario proves it).
# Requires+stop-propagation left lexa-api/telemetry permanently dead when a
# scenario cleanly stopped mosquitto. (V1RC FINDING A root fix, 2026-07-06.)
Wants=mosquitto.service

[Service]
# WS-1 sandboxing (V1.0 punch list, 2026-07-09): bounds blast radius on top
# of the existing User=lexa. StateDirectory grants /var/lib/lexa under
# ProtectSystem=strict for the scan-journal block configs/modbus.json ships
# (79ea417, unit 5.2) — without it the journal mkdir fails read-only at boot
# (bench round-1 finding, 2026-07-14). Same pattern as lexa-hub.service.
StateDirectory=lexa
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
RestrictAddressFamilies=AF_INET AF_UNIX
# Type=notify + WatchdogSec (TASK-008, review §11 item 5 / RSK-08): a wedged
# poll loop silently stops updating measurement freshness, which gates
# everything downstream (optimizer, compliance, telemetry). lexa-modbus
# sends READY=1 once reg.Start() (the poll loop) and publishMeasurements
# (the MQTT fan-out goroutine) are both running, and WATCHDOG=1 as the
# first statement of publishMeasurements' update-drain body
# (cmd/modbus/main.go) — covering registry, channel, and publish path in one
# kick site. A device poll error still emits an update each round (upd.Err),
# so unreachable devices alone do not starve the kick; only a wedged
# registry/channel does.
#
# WatchdogSec=60 sizing: >= 4x STOCK's ~10s poll x 3 devices (worst-case
# serialized poll round well under 60s in practice; parallel per-device poll
# keeps a single slow device from stalling the others). Kick wiring lives
# only in publishMeasurements — do not touch retryDevice/interlock code
# (radioactive per 05 §12) to add kicks elsewhere.
Type=notify
NotifyAccess=main
WatchdogSec=60
User=lexa
Group=lexa
ExecStart=/usr/local/sbin/lexa-modbus -config /etc/lexa/modbus.json
Restart=on-failure
RestartSec=5s

# Give Modbus connections time to settle before the next restart.

# journald rate cap (TASK-009, review §11 flash wear / RSK-14): normal
# steady-state logging is sparse (per-poll updates don't log unless
# upd.Err/REJECT-implausible/publish-error) but the hub re-commands battery
# setpoints every economic tick (FAST 3s), each logged on apply — ~10
# lines/30s per battery device baseline (estimated from the FAST tick cadence
# in CLAUDE.md/task background — pending on-bench measurement, see
# lexa-hub/docs/FLASH_BUDGET.md). Burst=300 gives >20x headroom for that plus
# a reconnect/poll-error storm (retryDevice logs on every drop/reconnect)
# without suppressing genuine fault-storm forensics.
LogRateLimitIntervalSec=30
LogRateLimitBurst=300

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lexa-modbus

[Install]
WantedBy=multi-user.target
