# systemd/lexa-api.service
[Unit]
Description=LEXA HTTP API (dashboard-facing /status + /logs)
Documentation=https://github.com/lexa/lexa-hub
After=mosquitto.service lexa-modbus.service lexa-northbound.service lexa-ocpp.service lexa-hub.service
Requires=mosquitto.service
Wants=lexa-modbus.service lexa-northbound.service lexa-ocpp.service lexa-hub.service

[Service]
Type=simple
User=lexa
Group=lexa
ExecStart=/usr/local/sbin/lexa-api -config /etc/lexa/api.json
Restart=on-failure
RestartSec=5s

StartLimitIntervalSec=60
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lexa-api

[Install]
WantedBy=multi-user.target
