# systemd/lexa-hub.service
[Unit]
# StartLimit* MUST live in [Unit] (systemd >=v230 ignores them in
# [Service]); the default 5-in-10s otherwise applies and a crash
# burst (e.g. the power-cut-retained scenario) permanently disables
# the service. (V1RC FINDING A, 2026-07-06.)
StartLimitIntervalSec=60
StartLimitBurst=5
Description=LEXA Hub Orchestrator (Energy Optimizer)
Documentation=https://github.com/lexa/lexa-hub
After=mosquitto.service lexa-modbus.service lexa-northbound.service
# Wants= not Requires=: a clean mosquitto stop/restart must NOT stop-propagate
# to the lexa services — they survive broker restarts via paho auto-reconnect
# (the crash-only/reconnect design; mqtt-broker-restart scenario proves it).
# Requires+stop-propagation left lexa-api/telemetry permanently dead when a
# scenario cleanly stopped mosquitto. (V1RC FINDING A root fix, 2026-07-06.)
Wants=mosquitto.service
Wants=lexa-modbus.service lexa-northbound.service lexa-ocpp.service

[Service]
# StateDirectory=lexa auto-creates /var/lib/lexa for the append-only
# journal (TASK-039/040). (V1RC FINDING D, 2026-07-06.)
StateDirectory=lexa
# WS-1 sandboxing (V1.0 punch list, 2026-07-09): bounds blast radius on top
# of the existing User=lexa. StateDirectory=lexa above already implies a
# read-write bind-mount for /var/lib/lexa (journal dir + breach snapshot,
# TASK-039/040/041) even under ProtectSystem=strict — systemd grants that
# automatically for RuntimeDirectory=/StateDirectory=/etc., no separate
# ReadWritePaths= needed. Nothing else in this process writes to disk.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
RestrictAddressFamilies=AF_INET AF_UNIX
# Type=notify + WatchdogSec (TASK-007, review §11 / §14 item 5): a crash
# already restarts the hub via Restart=on-failure below, but a live-but-
# wedged process (deadlocked tick loop, stalled paho publish, hung state
# read) previously never did. The hub sends READY=1 once its control loop
# is ticking (cmd/hub/main.go, right after eng.Start()) and WATCHDOG=1 from
# Engine.PlanObserver — which the orchestrator invokes ON THE ENGINE'S OWN
# CONTROL GOROUTINE every economic tick and on every safety pass that
# produces commands (internal/orchestrator/engine.go tick()/safetyTick()).
# That is deliberate: the kick rides the tick loop itself rather than an
# independent timer, so a wedge anywhere in the tick (ReadSystemState,
# Optimize, executePlan's synchronous MQTT publishes) starves the
# heartbeat instead of masking it.
#
# WatchdogSec sizing: >= 4x the slowest legitimate economic tick. The
# product ships STOCK (engine_interval_s=15 in configs/hub.json) but the
# bench QAs at FAST (3s) — size for STOCK since that's what's deployed to
# customers. Worst legitimate tick under a sick-but-alive broker: plan-log
# publish (5s PUBACK timeout) + a breach alert (5s) + up to 4 actuator
# publishes (4 x 5s) ~= 30s, still under the 60s budget below. 4 x 15s = 60s
# also lands exactly there, so WatchdogSec=60 is the floor that's safe for
# both timing modes while still catching a true wedge within about a
# minute. (This quantifies review §11's "synchronous publish waits" risk;
# TASK-046 removes that worst case properly — until then the watchdog just
# has to tolerate it.)
Type=notify
NotifyAccess=main
WatchdogSec=60
User=lexa
Group=lexa
ExecStart=/usr/local/sbin/lexa-hub -config /etc/lexa/hub.json
Restart=on-failure
RestartSec=5s


# journald rate cap (TASK-009, review §11 flash wear / RSK-14): FAST 3s tick
# logs ~2 lines/tick (plan summary + occasional decision line) ≈ 20 lines per
# 30s window (docs/refactor/tasks/TASK-009.md background estimate — pending
# on-bench measurement, see lexa-hub/docs/FLASH_BUDGET.md). Burst=400 is 20x
# that, generous enough to pass a fault-storm burst of breach/dedupe-reset
# logging uncapped while still catching a genuine runaway-logging bug
# (~100 lines/s would still trip it within a couple of intervals).
LogRateLimitIntervalSec=30
LogRateLimitBurst=400

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lexa-hub

[Install]
WantedBy=multi-user.target
