# systemd/lexa-telemetry.service
[Unit]
Description=LEXA CSIP Telemetry (MUP POST) Service
Documentation=https://github.com/lexa/lexa-hub
After=network-online.target mosquitto.service lexa-modbus.service
Requires=mosquitto.service
Wants=network-online.target lexa-modbus.service

[Service]
# Type=notify + WatchdogSec (TASK-008, review §11 item 5): lexa-telemetry has
# no tight control loop like northbound/modbus, so this is a weaker liveness
# signal by design (documented here and in cmd/telemetry/main.go): a 10s kick
# ticker lives as a case in the SAME select as the MUP post loop, so a
# wedged postMeasurements blocks the kick too rather than a free-running
# timer masking it. READY=1 is sent once MUP registration completes — each
# registerMUP call is a single bounded POST attempt (no internal retry loop),
# so this cannot hang startup on an unreachable server.
#
# WatchdogSec=60 sizing: >= 4x STOCK's 10s kick cadence with generous
# headroom; MUP posts happen every mup_post_rate_s (60s bench / 300s
# default), well outside the kick's own cadence since the kick ticker fires
# independently of the post ticker in the same select.
Type=notify
NotifyAccess=main
WatchdogSec=60
User=lexa
Group=lexa
ExecStart=/usr/local/sbin/lexa-telemetry -config /etc/lexa/telemetry.json
Restart=on-failure
RestartSec=10s

StartLimitIntervalSec=120
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lexa-telemetry

ReadOnlyPaths=/etc/lexa/certs

[Install]
WantedBy=multi-user.target
