# systemd/lexa-ocpp.service
[Unit]
# StartLimit* MUST live in [Unit] (systemd >=v230 ignores them in
# [Service]); the default 5-in-10s otherwise applies and a crash
# burst (e.g. the power-cut-retained scenario) permanently disables
# the service. (V1RC FINDING A, 2026-07-06.)
StartLimitIntervalSec=60
StartLimitBurst=5
Description=LEXA OCPP 2.0.1 CSMS (EV Charger Server)
Documentation=https://github.com/lexa/lexa-hub
After=network.target mosquitto.service
# Wants= not Requires=: a clean mosquitto stop/restart must NOT stop-propagate
# to the lexa services — they survive broker restarts via paho auto-reconnect
# (the crash-only/reconnect design; mqtt-broker-restart scenario proves it).
# Requires+stop-propagation left lexa-api/telemetry permanently dead when a
# scenario cleanly stopped mosquitto. (V1RC FINDING A root fix, 2026-07-06.)
Wants=mosquitto.service

[Service]
# WS-1 sandboxing (V1.0 punch list, 2026-07-09): bounds blast radius on top
# of the existing User=lexa. StateDirectory grants /var/lib/lexa under
# ProtectSystem=strict for the WP-13 pairing allowlist (allowlist_path in
# ocpp.json) — without it the first pairing approval fails read-only (same
# class as the bench round-1 lexa-modbus journal finding, 2026-07-14).
StateDirectory=lexa
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
RestrictAddressFamilies=AF_INET AF_UNIX
# Type=notify + WatchdogSec (TASK-008, review §11 item 5): lexa-ocpp's main
# goroutine only ever blocked on signal (srv.Start runs the WS listener in
# its own goroutine) — no natural liveness point exists. A 10s kick ticker
# gates on mc.IsConnected() (cmd/ocpp/main.go): liveness = process + MQTT
# connectivity. OCPP listener health itself is NOT probed here — a stuck WS
# accept loop with MQTT still up would NOT trip this watchdog (backlog:
# TASK-044 metrics probes listener health properly).
#
# WS-9.1: mc.IsConnected() alone is NOT enough — paho's AutoReconnect keeps
# it true for the client's entire reconnect retry loop, not just while
# actually connected, so a plain IsConnected() gate would never withhold a
# kick during a sustained outage. The kick is additionally gated on
# watchdog.DeafTracker (mqtt_deaf_restart_after_s in ocpp.json, default
# 300s): once the broker has been continuously unreachable that long, kicks
# stop, and WatchdogSec=60 restarts the service roughly
# mqtt_deaf_restart_after_s + 60s into the outage — accepted crash-only
# behavior (AD-011), not a bug; see PR for the mqtt-broker-restart /
# mqtt-broker-latency Mayhem reasoning.
Type=notify
NotifyAccess=main
WatchdogSec=60
User=lexa
Group=lexa
ExecStart=/usr/local/sbin/lexa-ocpp -config /etc/lexa/ocpp.json
Restart=on-failure
RestartSec=5s


# journald rate cap (TASK-009, review §11 flash wear / RSK-14): event-driven
# logging (StatusNotification/MeterValues/TransactionEvent), one line per
# OCPP message — an active charging session can push MeterValues every few
# seconds per station (estimated from evsim's sample cadence — pending
# on-bench measurement, see lexa-hub/docs/FLASH_BUDGET.md). Burst=300 gives
# generous headroom over a handful of concurrently-charging stations without
# suppressing session forensics.
LogRateLimitIntervalSec=30
LogRateLimitBurst=300

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lexa-ocpp

# OCPP listens on a privileged port only when port < 1024.
# For port 8887 (default) no special capability is needed.
# AmbientCapabilities=CAP_NET_BIND_SERVICE  # uncomment if port < 1024

[Install]
WantedBy=multi-user.target
